# docker-backend/gdbe-supervisor.conf
[program:gdbe-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/docker/code/uwsgi.ini

[program:gdbe-nginx]
command = /usr/sbin/nginx
